# scripts/cloudflared-quick-tunnel.service
[Unit]
Description=Cloudflare Quick Tunnel (ResumePilot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@RUN_USER@
ExecStart=/usr/bin/cloudflared tunnel --url http://127.0.0.1:@APP_PORT@
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
